# systemd/user/meshcore-mqtt-restart.timer
[Unit]
Description=Trigger daily restart of MeshCore MQTT Bridge

[Timer]
OnCalendar=*-*-* 03:30:00
Persistent=true
RandomizedDelaySec=10m
Unit=meshcore-mqtt-restart.service

[Install]
WantedBy=timers.target
