# systemd/user/meshcore-mqtt-restart.service
[Unit]
Description=Daily restart of MeshCore MQTT Bridge

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user restart meshcore-mqtt.service
